# collections/ansible_collections/ejsmit/containers/roles/docker-nextcloud/templates/nextcloud_cron.service


[Unit]
Description=Run Nextcloud cron job every 5 mins

[Service]
ExecStart=docker exec --user www-data nextcloud php -f /var/www/html/cron.php
